# standalone_web_service/etc/rc.local
#!/bin/sh
/bin/bash /ritar-bms/run.sh &
exit 0;
